# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true   # Run in the foreground (not as a daemon)

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:mysql]
command=/usr/bin/mysqld_safe
